# brightness.cron

#* * * * * root source /root/.bash_profile; /opt/train-alerts/bin/brightness.sh >> /var/log/cron.log 2>&1
* * * * * root /opt/train-alerts/bin/brightness.sh >> /var/log/cron.log 2>&1
